# radiod.service
# Radio systemd script
# $Id: radiod.service,v 1.5 2025/06/26 10:24:20 bob Exp $
[Unit]
Description=Radio daemon
After=network.target

[Service]
Type=simple
ExecStart=/usr/share/radio/radiod.py nodaemon
PIDFile=/var/run/radiod.pid
#Restart=on-failure
#RestartSec=5s

[Install]
WantedBy=multi-user.target
